# etc/systemd/system/capbot-sp500_5m_rth.service
[Unit]
Description=CapBot sp500_5m_rth
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/capbot_full

# Optional: systemd env file (can be empty). Keep secrets out of the repo.
EnvironmentFile=-/home/ubuntu/sp500_5m_rth.secrets

# Real-time logs
Environment=PYTHONUNBUFFERED=1

# Avoid stale lock from previous crash
ExecStartPre=/bin/rm -f /home/ubuntu/.capbot_lock_sp500_5m_rth.lock

ExecStart=/opt/capbot_full/.venv/bin/python -m capbot.app.main_sp500_rth run --config /opt/capbot_full/configs/sp500_5m_rth.json --secrets /home/ubuntu/sp500_5m_rth.secrets

Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
